# packer-ami/webapplication.service
[Unit]
Description=<Flask web application>
After=network.target

[Service]
User=ubuntu
ExecStart=python /Application/main.py
Restart=always